# systemd-startup/amanda-client@.service.in
[Unit]
Description=Amanda Backup Client Connection
ConditionPathExists=/etc/amanda-security.conf
Requires=amanda-client.socket

[Service]
User=@CLIENT_LOGIN@
Group=@SETUID_GROUP@
# no auto-restart possible...
ExecStart=@amlibexecdir@/amandad -auth=bsdtcp amdump senddiscover
StandardInput=socket
