# systemd-startup/amanda-client.socket.in
[Unit]
Description=Amanda Backup Client Socket
Requires=network-online.target

[Socket]
# only for the Unix sockets
SocketUser=@CLIENT_LOGIN@
SocketGroup=@SETUID_GROUP@
SocketMode=0600
ListenStream=10080
ListenStream=@amdatadir@/amandad.sock
Accept=yes
KeepAlive=true
RemoveOnStop=yes
#
# dont wait two hours
#
KeepAliveTimeSec=120
# KeepAliveIntervalSec=75
# KeepAliveProbes=9
#
# speed up unreplied bursts
#
NoDelay=true
#
# set TOS for TCP outputs if needed (0-6 only normally allowed)
#
#Priority= 
#
# keep the daemon asleep until real data arrives
#
DeferAcceptSec=5
# ReceiveBuffer=
# SendBuffer=
IPTOS=throughput
ReusePort=false
# SELinuxContextFromNet=
# TCPCongestion=bbr
MaxConnections=3

[Install]
# started at boot time
WantedBy=sockets.target
