# systemd-startup/amanda-server@.service.in
[Unit]
Description=Amanda Backup Service Connection
ConditionPathExists=/etc/amanda-security.conf
Requires=amanda-server.socket

[Service]
User=@CLIENT_LOGIN@
Group=@SETUID_GROUP@
# no restart possible
PermissionsStartOnly=yes
ExecStart=@amlibexecdir@/amandad -auth=bsdtcp amdump senddiscover amindexd amidxtaped
StandardInput=socket
