# example/systemd/amanda-rest-server-pxy.socket.in
[Unit]
Description=Amanda Backup Server Socket
# sends starts or restarts to
BindsTo=amanda-rest-server.service
# delayed if no network active
After=network-online.target

[Socket]
# local only!
ListenStream=[::1]:5000
##SELinuxContext=system_u:system_r:amanda_data_t:s0
Accept=no
MaxConnections=5

[Install]
# receives starts from ..
RequiredBy=amanda-rest-server.service
