# example/systemd/amanda-lserver@.service.in
[Unit]
Description=Amanda Backup Service Connection
Requires=amanda-lserver.socket

[Service]
User=@CLIENT_LOGIN@
Group=@SETUID_GROUP@
ExecStart=@amlibexecdir@/amandad -auth=local amdump senddiscover amindexd amidxtaped
StandardInput=socket
