# example/systemd/amanda-rest-server-pxy.service.in
[Unit]
Description=Amanda Rest Server Proxy
# just fail if missing...
ConditionPathExists=@amdatadir@/amanda-rest.sock
# starts up, waits for but can be stopped by ...
BindsTo=__sub-amanda-rest-server-pxy.socket
BindsTo=__sub-amanda-rest-server.service
After=network-online.target
After=__sub-amanda-rest-server.service
After=__sub-amanda-rest-server-pxy.socket

[Service]
User=@CLIENT_LOGIN@
Group=@SETUID_GROUP@
ExecStart=/usr/lib/systemd/systemd-socket-proxyd @amdatadir@/amanda-rest.sock
PrivateTmp=yes

[Install]
Alias=amanda-rest-server-pxy.service
