# example/systemd/amanda-server.socket.in
[Unit]
Description=Amanda Backup Server Socket
# delayed if no network active
After=network-online.target

[Socket]
ListenStream=10080
Accept=yes
KeepAlive=true
#
# dont wait two hours
#
KeepAliveTimeSec=120
#KeepAliveIntervalSec=75
#KeepAliveProbes=9
#
# speed up unreplied bursts
#
NoDelay=true
#
# set TOS for TCP outputs if needed (0-6 only normally allowed)
#
#Priority= 
#
# keep the daemon asleep until real data arrives
#
DeferAcceptSec=5
#ReceiveBuffer=
#SendBuffer=
IPTOS=throughput
ReusePort=false
#SELinuxContextFromNet=
#TCPCongestion=bbr
MaxConnections=10

[Install]
Alias=amanda-server.socket
WantedBy=sockets.target
