# example/systemd/amanda-client.socket.in
[Unit]
Description=Amanda Backup Client Socket

[Socket]
ListenStream=10080
Accept=yes
KeepAlive=true
######## dont wait two hours
KeepAliveTimeSec=120
#KeepAliveIntervalSec=75
#KeepAliveProbes=9
######## speed up unreplied bursts
NoDelay=true
######## set TOS for TCP outputs if needed (0-6 only normally allowed)
#Priority= 
######## don't fork until real data arrives
DeferAcceptSec=5
IPTOS=throughput
# allow localhost-only socket elsewhere
ReusePort=true
# SELinuxContextFromNet=
# TCPCongestion=bbr
MaxConnections=10

[Install]
# started at boot time
WantedBy=sockets.target
