# example/systemd/amanda-lserver.socket.in
[Unit]
Description=Amanda Backup Server Localhost-Socket

[Socket]
ListenStream=[::1]:10080
Accept=yes
ReusePort=true
MaxConnections=100

[Install]
# started at boot time
WantedBy=sockets.target
