# example/systemd/amanda-server@.service.in
[Unit]
Description=Amanda Backup Service Connection
ConditionPathExists=@DEFAULT_SECURITY_FILE@

[Service]
User=@CLIENT_LOGIN@
Group=@SETUID_GROUP@
# no restart possible
ExecStart=@amlibexecdir@/amandad -auth=bsdtcp amdump senddiscover amindexd amidxtaped
StandardInput=socket
