# example/systemd/amanda-client.socket.in
[Unit]
Description=Amanda Backup Client Socket

[Socket]
ListenStream=10080
Accept=yes
RemoveOnStop=yes
ReusePort=true
MaxConnections=10
##SELinuxContext=system_u:system_r:amanda_data_t:s0
KeepAlive=true
######## dont wait two hours
KeepAliveTimeSec=120
######## speed up unreplied bursts
NoDelay=true
######## set TOS for TCP outputs if needed (0-6 only normally allowed)
#Priority= 
######## keep the daemon asleep until real data arrives
DeferAcceptSec=5
IPTOS=throughput
# allow localhost-only socket elsewhere

[Install]
WantedBy=multi-user.target
