# example/systemd/amanda-lserver.socket.in
[Unit]
Description=Amanda Backup Server Localhost-Socket

[Socket]
ListenStream=127.0.0.1:10080
Accept=yes
RemoveOnStop=yes
ReusePort=true
MaxConnections=100
##SELinuxContext=system_u:system_r:amanda_data_t:s0

[Install]
WantedBy=multi-user.target
