# example/systemd/amanda-rest-server.service.in
[Unit]
Description=Amanda Backup Rest Service

[Service]
Type=forking
User=@CLIENT_LOGIN@
Group=@SETUID_GROUP@
WorkingDirectory=@amdatadir@
PIDFile=@amdatadir@/tmp/rest-api-pid
# restart as needed
RestartSec=10
Restart=always

# need a delay to detect real start time
ExecStart=/usr/sbin/amanda-rest-server start
ExecStartPost=/bin/bash -xvc "for i in 1 2 3 4 5; do [ -s @amdatadir@/tmp/rest-api-pid ] && break; sleep 1; done"

ExecStop=/usr/sbin/amanda-rest-server stop
PrivateTmp=true

[Install]
# start only if amanda sockets are running too
WantedBy=__sub-amanda-server.socket
WantedBy=__sub-amanda-lserver.socket
